# main/firewall/debian/lucid/zentyal-firewall.prerm
#!/bin/bash

case "$1" in
	remove)
        invoke-rc.d zentyal firewall stop
		;;
esac

exit 0
