# main/firewall/debian/lucid/zentyal-firewall.prerm
#!/bin/bash

set -e

case "$1" in
	remove)
        invoke-rc.d zentyal firewall stop || true
		;;
esac 

#DEBHELPER#
