# extra/openchange/debian-sope/libapache2-mod-ngobjweb.postinst
#!/bin/bash

set -e

if [ -e /etc/apache2/httpd.conf ]; then
	a2enmod ngobjweb
fi
if [ -x /usr/sbin/invoke-rc.d ] && [ -x /etc/init.d/apache2 ]; then
	invoke-rc.d apache2 force-reload || true
elif [ -x /etc/init.d/apache2 ]; then
	/etc/init.d/apache2 force-reload || true
fi

#DEBHELPER#
