# extra/installer/data/rc.local
#!/bin/sh

/etc/init.d/slim restart

exit 0
