# main/core/debian/precise/zentyal.cleantmp.upstart
description "Clean Zentyal temp directory"

start on filesystem

task

exec /bin/rm -rf /var/lib/zentyal/tmp/*
