# client/ipsec/debian/ebox-ipsec.postrm
#!/bin/bash

case "$1" in
	remove)
		invoke-rc.d ipsec stop || true

        dpkg-trigger --no-await ebox
		;;
esac

#DEBHELPER#
