# main/core/extra/files.list
keyring.gpg /usr/share/zentyal
ebox.passwd /var/lib/zentyal/conf
cron/zentyal /etc/cron.daily
cron/zentyal-events /etc/cron.hourly
cron/90zentyal-manage-logs /etc/cron.hourly
css /usr/share/zentyal/stubs
