# main/core/extra/files.list
keyring.gpg /usr/share/zentyal
zentyal-qa.asc /etc/apt/trusted.gpg.d
99zentyal /etc/apt/apt.conf.d
ebox.passwd /var/lib/zentyal/conf
cron/zentyal /etc/cron.daily
cron/90zentyal-manage-logs /etc/cron.hourly
css /usr/share/zentyal/stubs
mysql/zentyal.cnf /etc/mysql/conf.d
eboxlog.conf /var/lib/zentyal/conf
zs /usr/bin
bash-completion/zs /usr/share/bash-completion/completions
