# client/openvpn/debian/ebox-openvpn.postrm
#!/bin/bash

#DEBHELPER#

case "$1" in
	remove)

		kill  `pidof gconfd-2` >/dev/null 2>&1 || true
	
		# openvpn changes network interfaces so we must restart network module to propagate the changes
        invoke-rc.d ebox network restart   || true 

        # Restart eBox apache
        invoke-rc.d ebox apache restart || true

		;;		
esac


