# client/ebox/debian/ebox.prerm
#!/bin/bash

set -e

invoke-rc.d ebox stop || true

#DEBHELPER#

