# client/ebox/debian/ebox.prerm
#!/bin/bash

if [ "$1" = remove ] || [ "$1" = purge ]; then
    invoke-rc.d ebox stop || true
fi

#DEBHELPER#

