# client/ebox/debian/ebox.cron.hourly
#!/bin/bash

invoke-rc.d ebox events restart
