# main/openchange/debian/precise/zentyal-openchange.postinst
#!/bin/bash

case "$1" in
    configure)

        /usr/share/zentyal/initial-setup --no-restart openchange $2

        # restart module
        invoke-rc.d zentyal openchange restart || true

        dpkg-trigger --no-await zentyal-core
    ;;
esac

exit 0
