# client/remoteservices/debian/lucid/ebox-remoteservices.prerm
#!/bin/bash

#DEBHELPER#

if [ "$1" = remove ] || [ "$1" = purge ]; then
    invoke-rc.d ebox remoteservices stop
    /usr/share/ebox/ebox-unsubscribe
fi

exit 0
