# client/usersandgroups/tools/ebox-ad-sync.cron
# FIXME: Change this to 5 minutes instead of 1
*/1 * * * * /usr/share/ebox-usersandgroups/ebox-ad-sync

