# client/usersandgroups/tools/ebox-ad-sync.cron
*/5 * * * * root /usr/share/ebox-usersandgroups/ebox-ad-sync
