# main/remoteservices/debian/lucid/zentyal-remoteservices.prerm
#!/bin/bash

case "$1" in
    remove)
        invoke-rc.d zentyal remoteservices stop
        /usr/share/zentyal-remoteservices/unsubscribe
    ;;
esac

exit 0
