# main/users/src/scripts/ad-sync.cron
*/5 * * * * root /usr/share/zentyal-users/ad-sync
