# extra/installer/zenbuntu-desktop/rc.local
#!/bin/sh

initctl emit zentyal-lxdm

exit 0
