# main/mail/extra/files.list
dovecot-pam /usr/share/zentyal-mail
zentyal-autoexpunge /etc/cron.daily/zentyal-autoexpunge
