# client/mail/debian/ebox-mail.postrm
#!/bin/bash

case "$1" in
	remove)
		kill `pidof gconfd-2` >/dev/null 2>&1 || true
		invoke-rc.d postfix stop || true
		;;		
esac

#DEBHELPER#
