# client/dns/debian/ebox-dns.postrm
#!/bin/bash

case "$1" in
	remove)
	        /usr/bin/ebox-clean-gconf dns || true
		kill `pidof gconfd-2` >/dev/null 2>&1 || true
		invoke-rc.d bind9 stop || true
		;;		
esac

#DEBHELPER#
