# client/dns/debian/ebox-dns.postinst
#!/bin/bash

set -e

#DEBHELPER#



case "$1" in
	configure)
		kill `pidof gconfd-2` >/dev/null 2>&1 || true

        # Migrate data if needed
        /usr/share/ebox/ebox-migrate /usr/share/ebox-dns/migration/

		invoke-rc.d ebox dns restart

        # Restart eBox apache
        invoke-rc.d ebox apache restart || true

esac



exit 0
