# main/webserver/debian/lucid/zentyal-webserver.postinst
#!/bin/bash

#DEBHELPER#

case "$1" in
	configure)

		# Migrate data if needed
		/usr/share/zentyal/ebox-migrate /usr/share/zentyal-webserver/migration/

		invoke-rc.d zentyal webserver restart || true

        dpkg-trigger --no-await zentyal-core
esac



exit 0
