# client/ebox/debian/lucid/ebox.cron.daily
#!/bin/sh

# Check latest version
URL="http://update.ebox-platform.com/last-1.5.txt"
TIMEOUT="300"
DEST="/var/lib/ebox/latestversion"
wget --timeout $TIMEOUT $URL -O $DEST >/dev/null 2>&1

# Gather reporting info
/usr/share/ebox/ebox-gather-reportinfo
