# client/ebox/debian/lucid/ebox.cron.hourly
#!/bin/bash

invoke-rc.d ebox events restart > /dev/null 2>&1
