# client/software/debian/ebox-software.postinst
#!/bin/bash

#DEBHELPER#
case "$1" in
    configure)
        kill `pidof gconfd-2` >/dev/null 2>&1 || true

        # migrate data if needed
        /usr/share/ebox/ebox-migrate /usr/share/ebox-software/migration/	

        invoke-rc.d ebox apache restart || true
    ;;

esac



exit 0
