# main/squid/debian/precise/ebox.dansguardian.upstart
pre-start script
    invoke-rc.d dansguardian stop || true
end script

start on runlevel [2345]
stop on runlevel [!2345]
manual
respawn

exec /usr/sbin/dansguardian -N
