# main/mailfilter/debian/zentyal.learnspamd.service
[Unit]
Description=Zentyal learn spam daemon

[Service]
ExecStart=/usr/share/zentyal-mailfilter/learnspamd
Restart=on-failure
