# extra/installer/zenbuntu-desktop/debian/postinst
#!/bin/bash

ln -sf /usr/share/plymouth/themes/zentyal/zentyal.plymouth /etc/alternatives/default.plymouth
ln -sf /usr/share/plymouth/themes/zentyal/zentyal.grub /etc/alternatives/default.plymouth.grub

if [ "$1" == "configure" ] && [ -z "$2" ]
then
    systemctl disable lxdm.service || true
fi
