# main/dns/extra/files.list
fix-744304.conf /etc/systemd/system/named.service.wants