# main/core/debian/zentyal-core.prerm
#!/bin/bash

set -e

#DEBHELPER#

# FIXME: check if this works properly
if [ "$1" = remove ] || [ "$1" = purge ]; then
    invoke-rc.d zentyal stop || true
fi

exit 0
