# main/core/debian/zentyal.loggerd.service
[Unit]
Description=Zentyal logger daemon

[Service]
ExecStart=/usr/share/zentyal/loggerd
Restart=on-failure
