# client/dhcp/debian/ebox-dhcp.postrm
#!/bin/bash

#DEBHELPER#

case "$1" in
	remove)
		kill `pidof gconfd-2` >/dev/null 2>&1 || true
		invoke-rc.d ebox logs restart || true

        # Restart eBox apache
        invoke-rc.d ebox apache restart || true

		;;		
	
	purge)
		# drop log table
		/usr/share/ebox/ebox-sql-table drop leases || true

		;;

esac
