# client/jabber/debian/ebox-jabber.postrm
#!/bin/bash

case "$1" in
	remove)
		kill `pidof gconfd-2` >/dev/null 2>&1 || true
		invoke-rc.d  jabberd2-ldap-bdb stop || true

                # Generate configuration
                /usr/lib/ebox-usersandgroups/ebox-init-ldap genconfig

                # Restart slapd
                invoke-rc.d slapd restart

		/usr/bin/ebox-clean-gconf jabber || true
		;;		
esac

#DEBHELPER#
