# main/network/debian/lucid/zentyal-network.postinst
#!/bin/bash

set -e

#DEBHELPER#

case "$1" in
    configure)
        # inital setup
        /usr/share/zentyal/initial-setup network $2

        # restart module
        invoke-rc.d zentyal network restart || true

        dpkg-trigger --no-await zentyal-core
    ;;
esac

exit 0
