# main/core/debian/lucid/ebox.cron.hourly
#!/bin/bash

invoke-rc.d zentyal events restart > /dev/null 2>&1
