# main/core/debian/lucid/ebox.cron.daily
#!/bin/sh

# Check latest version
URL="http://update.zentyal.org/last-2.0.txt"
TIMEOUT="300"
DEST="/var/lib/ebox/latestversion"
wget --timeout $TIMEOUT $URL -O $DEST >/dev/null 2>&1

# Gather reporting info
/usr/share/zentyal/ebox-gather-reportinfo
# Consolidate this reporting info
/usr/share/zentyal/ebox-consolidate-reportinfo
