# main/core/debian/lucid/zentyal-core.prerm
#!/bin/bash

set -e

if [ "$1" = remove ] || [ "$1" = purge ]; then
    invoke-rc.d zentyal stop || true
fi
