# main/users/debian/lucid/zentyal-users.postinst
#!/bin/bash

set -e
#DEBHELPER#

case "$1" in
    configure)

        /usr/share/zentyal/initial-setup users $2

        # Restart slapd after the the configuration has been
        # generated by the module initial setup
        invoke-rc.d slapd restart || true

        dpkg-trigger --no-await zentyal-core
    ;;
esac

exit 0
