# main/users/debian/lucid/zentyal-users.postinst
#!/bin/bash

set -e

#DEBHELPER#

case "$1" in
    configure)
        # initial setup
        /usr/share/zentyal/initial-setup users $2

        # restart slapd after the the configuration has been
        # generated by the module initial setup
        #invoke-rc.d slapd restart || true #XXX ?

        # restart module
        invoke-rc.d zentyal users restart || true

        dpkg-trigger --no-await zentyal-core
    ;;
esac

exit 0
