# client/dns/debian/lucid/ebox-dns.postinst
#!/bin/bash

set -e

#DEBHELPER#



case "$1" in
    configure)
    # Migrate data if needed
    /usr/share/ebox/ebox-migrate /usr/share/ebox-dns/migration/

    invoke-rc.d ebox dns restart

    dpkg-trigger --no-await ebox
esac



exit 0
