# client/samba/debian/ebox-samba.postrm
#!/bin/bash

case "$1" in
	remove)
		kill `pidof gconfd-2` >/dev/null 2>&1 || true
		invoke-rc.d samba stop || true
		
		# Generate configuration
		/usr/lib/ebox-usersandgroups/ebox-init-ldap genconfig

		# Restart slapd
		invoke-rc.d slapd restart

		;;		
	purge)
	 	test -d /home/samba && rm -rf /home/samba

	        ;;    
esac

#DEBHELPER#
