# main/webmail/debian/lucid/zentyal-webmail.postinst
#!/bin/bash

#DEBHELPER#

case "$1" in
	configure)

        # migrate data
#       uncoment next line when is needed to migrate data
#       /usr/share/ebox/ebox-migrate /usr/share/ebox-webmail/migration

        # restart module
        invoke-rc.d ebox webmail restart

        dpkg-trigger --no-await zentyal-core
esac

exit 0
