# main/openvpn/debian/lucid/zentyal-openvpn.postrm
#!/bin/bash

case "$1" in
	remove)
        # openvpn changes network interfaces so we must restart network module to propagate the changes
        invoke-rc.d zentyal network restart

        dpkg-trigger --no-await zentyal-core
        ;;
    purge)
        /usr/share/zentyal/purge-module openvpn
        ;;
esac

exit 0
