# main/samba/extra/files.list
krb5Key.asn /usr/share/zentyal-samba
cron/zentyal-samba /etc/cron.d
