# client/firewall/debian/ebox-firewall.prerm
#!/bin/bash

#DEBHELPER#

case "$1" in
	remove)
        invoke-rc.d ebox firewall stop || true
		;;
esac 
