# client/trafficshaping/debian/natty/ebox-trafficshaping.postinst
#!/bin/bash

set -e

#DEBHELPER#

case "$1" in
    configure)
        invoke-rc.d ebox trafficshaping restart || true

        dpkg-trigger --no-await ebox
    ;;
esac

exit 0
