# main/users/debian/zentyal-users.postinst
#!/bin/bash

set -e

#DEBHELPER#

case "$1" in
    configure)
        # initial setup
        /usr/share/zentyal/initial-setup --no-restart users $2

        # restart module
        invoke-rc.d zentyal users restart || true

        dpkg-trigger --no-await zentyal-core
    ;;
esac

exit 0
