# contrib/egroupware/debian/lucid/ebox-egroupware.postinst
#!/bin/bash

#DEBHELPER#

case "$1" in
    configure)
        # Migrate data if needed
        /usr/share/ebox/ebox-migrate /usr/share/ebox-egroupware/migration/

        # Restart webserver module
        invoke-rc.d ebox webserver restart || true

        dpkg-trigger --no-await ebox
    ;;
esac

exit 0
