# client/ebox/debian/natty/ebox.cron.daily
#!/bin/sh

# check latest version
URL="http://update.zentyal.org/last-2.0.txt"
TIMEOUT="300"
DEST="/var/lib/ebox/latestversion"
wget --timeout $TIMEOUT $URL -O $DEST >/dev/null 2>&1

# gather reporting info
/usr/share/ebox/ebox-gather-reportinfo

# consolidate this reporting info
/usr/share/ebox/ebox-consolidate-reportinfo
