# main/samba/debian/precise/zentyal.nmbd.upstart
description "NetBIOS name server"
author      "Zentyal Packaging Maintainers <pkg-team@zentyal.com>"

start on started samba4
stop on stopping samba4

expect fork
normal exit 0

exec nmbd -D
