# main/samba/debian/precise/zentyal.s4sync.upstart
start on started samba4
stop on stopping samba4

exec /usr/share/zentyal-samba/s4sync
respawn
